# paperboy-project-generator/generators/frontend/templates/frontend/supervisord.conf
[supervisord]
nodaemon = true
user = node

[program:paperboy]
command = node_modules/.bin/paperboy -c paperboy/paperboy.docker.json
directory = /home/node/app
autostart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/node",USER="node"
autorestart=true

[program:frontend]
command = npm start
directory = /home/node/app
autostart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/node",USER="node"
autorestart=true